# services/piper/systemd/piper.service
[Unit]
Description=Aila Piper TTS Service
After=network.target sound.target

[Service]
Type=simple
User=aila
Group=aila
EnvironmentFile=/etc/aila/env.d/piper.conf
WorkingDirectory=/opt/aila/piper
ExecStart=/usr/bin/python3 -m main.server --voice /etc/aila-piper/voice.json
Restart=on-failure
RestartSec=3
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
