# services/coqui/systemd/coqui.service
[Unit]
Description=Coqui TTS Service
After=network.target sound.target

[Service]
Type=simple
User=aila
Group=aila
EnvironmentFile=/etc/aila/env.d/coqui.conf
WorkingDirectory=/opt/aila/coqui
ExecStart=/opt/aila/coqui/venv/bin/python -m server --config /etc/aila-coqui/config.yaml
Restart=on-failure
RestartSec=3
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
